# slavefw.service
[Unit]
Description=Slave FW Autostarter

[Service]
Type=oneshot
ExecStart=/root/slavefw.sh

[Install]
WantedBy=multi-user.target